# vm/quadlets/hermes-proxy.container
[Unit]
Description=Hermes security proxy
After=hermes-network.service
Wants=hermes-network.service

[Container]
Image=ghcr.io/{{ ghcr_owner }}/hermes-proxy:latest
ContainerName=hermes-proxy
PublishPort=8000:8000
EnvironmentFile=/etc/containers/systemd/hermes.env
Network=hermes.network
Volume=hermes-proxy-workspace:/app/workspace
Volume=hermes-proxy-data:/app/data
HealthCmd=["python3", "-c", "import urllib.request; urllib.request.urlopen('http://localhost:8000/health')"]
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=10s
AutoUpdate=registry

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
