# vm/quadlets/hermes-discord.container
[Unit]
Description=Hermes Discord bot
After=hermes-proxy.service
Wants=hermes-proxy.service

[Container]
Image=ghcr.io/{{ ghcr_owner }}/hermes-discord:latest
ContainerName=hermes-discord
EnvironmentFile=/etc/containers/systemd/hermes.env
Environment=PROXY_URL=http://hermes-proxy:8000
Network=hermes.network
HealthCmd=["python3", "-c", "import urllib.request; urllib.request.urlopen('http://localhost:8001/health')"]
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=15s
AutoUpdate=registry

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
